# deploy/systemd/youth-together-pipeline.service
[Unit]
Description=Youth Together content pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=youthtogether
Group=youthtogether
WorkingDirectory=/opt/youth-together
Environment=TZ=Asia/Seoul
EnvironmentFile=-/etc/youth-together.env
ExecStart=/opt/youth-together/deploy/linux/run_pipeline.sh
TimeoutStartSec=1800
Nice=10
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
